# deploy/ai-support/systemd/denstock-ai-launcher@.service
[Unit]
Description=DenisStock AI launcher request
Requires=denstock-ai-firewall.service denstock-ai-proxy.service
After=denstock-ai-firewall.service denstock-ai-proxy.service

[Service]
Type=exec
User=root
Group=root
UMask=0077
ExecStart=/usr/local/sbin/denstock-ai-launcher socket-serve
StandardInput=socket
StandardOutput=socket
StandardError=null
TimeoutStartSec=330s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
IPAddressDeny=any
IPAddressAllow=localhost
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_KILL CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_NET_ADMIN
AmbientCapabilities=
ReadOnlyPaths=/etc/denstock-ai/launcher.json /usr/local/bin/codex
ReadWritePaths=/var/lib/denstock-ai/codex-home /var/lib/denstock-ai/requests /run/denstock-ai/locks
InaccessiblePaths=-/opt/denstock -/opt/backups -/var/backups -/var/lib/docker -/var/lib/postgresql -/run/docker.sock -/run/postgresql -/var/log/denstock-backup.log
